# nodist/etc/cron.d/crontab-sysadmin
*/2 * * * *     root    /usr/local/sbin/proc.sh
@daily  root            /usr/local/sbin/backup_mysql.sh localhost
@hourly root            /usr/bin/find /var/log/ps/ -maxdepth 1 -type d -mtime +30 -exec rm -rf "{}" \;
@hourly root            echo `date +\%F` > /var/run/backup_timestamp
