# nodist/etc/cron.d/crontab-sysadmin-backup-timestamp
@hourly root echo `date +\%F` > /var/run/backup_timestamp
